# refreshLab.service
[Unit]
Description="Hello World script"

[Service]
ExecStart=/usr/local/bin/refreshLab.sh